# etc/systemd/quirck.socket
[Unit]
Description=quirck socket

[Socket]
ListenStream=/run/quirck.sock
SocketUser=www-data
SocketMode=600

[Install]
WantedBy=sockets.target
