# etc/systemd/quirck-web.socket
[Unit]
Description=quirck socket

[Socket]
ListenStream=/run/quirck/quirck.sock
SocketMode=666

[Install]
WantedBy=sockets.target
